# fluent-bit/lib/monkey/monkey.service
[Unit]
Description=Monkey HTTP Server
Requires=network.target
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/sbin/monkey --daemon
PIDFile=/
Restart=always

[Install]
WantedBy=multi-user.target
